# confluence/systemd/postgres.service
[Unit]
Description=%p Service
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=/usr/bin/docker pull %p:9.6
ExecStartPre=-/usr/bin/docker rm --force %p
ExecStart=/usr/bin/docker run --name %p \
  -e POSTGRES_PASSWORD=postgres \
  --net host \
  -v /opt/atlassian/postgres:/var/lib/postgresql/data \
  %p:9.6
ExecStop=-/usr/bin/docker stop --time 3 %p

[Install]
WantedBy=multi-user.target
